# kits/hermes-webui/src/configs/com.hermes.caddy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermes.caddy</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/caddy</string>
        <string>run</string>
        <string>--config</string>
        <string>/Users/dodhya/.hermes/caddy/Caddyfile</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/dodhya</string>
        <key>XDG_DATA_HOME</key>
        <string>/Users/dodhya/.hermes/caddy/data</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/dodhya/.hermes/caddy/caddy.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/dodhya/.hermes/caddy/caddy.log</string>
</dict>
</plist>
